# services/nginx/run
#!/bin/bash
set -euo pipefail
exec nginx -c nginx.conf -p $(pwd)
