# services/web/run
#!/bin/bash
set -euxo pipefail
exec 2>&1

cd /opt/ocfweb
exec /opt/ocfweb/venv/bin/gunicorn \
    -b 0.0.0.0:8000 \
    -w 4 \
    ocfweb.wsgi
